# etc/profile.d/EIGENSOFT.sh
#EIG5.0.2

if ! echo ${PATH} | grep -q /usr/local/EIGENSOFT/bin ; then
        export PATH=${PATH}:/usr/local/EIGENSOFT/bin
fi

