# etc/profile.d/pique.sh
#@(#)pique.sh  2015-11-18  A.J.Travis and A.Douglas

#
# Environment for running pique
#

export PIQUE_HOME=/usr/local/pique
export PATH=$PIQUE_HOME/bin:$PATH
